# deploy/systemd/he-uc-lending.service
[Unit]
Description=HE UC Lending encrypted job receiver
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=/opt/he_uc_lending
EnvironmentFile=/etc/he-uc-lending.env
ExecStart=/usr/bin/python3 /opt/he_uc_lending/code/server/web/he_job_server.py --host ${HE_WEB_HOST} --port ${HE_WEB_PORT} --build-dir /opt/he_uc_lending/build --jobs-dir /opt/he_uc_lending/server_jobs/web
Restart=always
RestartSec=3
TimeoutStopSec=20
KillSignal=SIGINT
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
